# install/templates/myclash.service
[Unit]
Description=MyClashShell MCS manager (systemd → Python; mihomo as mcs/bin/clash via subprocess)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
TimeoutStartSec=30
WorkingDirectory=<|PLACEHOLD_1|>
Environment=MYCLASH_ROOT_PWD=<|PLACEHOLD_2|>
ExecStart=<|PLACEHOLD_1|>/venv/bin/python3 <|PLACEHOLD_2|>/scripts/runtime/mcs_manager.py
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10s

[Install]
# 用户级 systemd（systemctl --user enable）使用 default.target，勿用 multi-user.target
WantedBy=default.target
